# deb-meta-package/postinst
#!/bin/sh 

systemctl enable prex
systemctl enable linkbot-hub-manager
/etc/cron.daily/auto-update-pylinkbot

# Add the new Barobo public key

echo "-----BEGIN PGP PUBLIC KEY BLOCK-----
Version: GnuPG v1

mQENBFfAkpYBCACogXlFdy8TeB9yNJKRA7ctHALhNWmE5FDDihgtYgY9h4VvGkoD
mhqPxYoJrIjhZkb4QGSCP0x8VXa0aFpW5kMwGWc+Q8Yh7dVbczarSivuYHLoHZBK
yaTz7jLD47D8zhPOtcP7znUw3enWGjvT7wN/vTI2i3g4v4GmEGy9gTO36Nsu/cjl
0T1089WNNZdjE2MALG9YadB3fe/V86DgadpjYf9OmsB9t+PiIMA/G/0Tb7gc2z6C
nS7JPbfybwO1bbSgfFVXbDPxsjiSQGFhC8MygCv0xxvNYj4uytiX2aVCs4Vyr03G
8erQTWGJ+DwBbkD+c06+dhYdwaTC03pzXGzVABEBAAG0HUJhcm9ibyBJbmMuIDxp
bmZvQGJhcm9iby5jb20+iQE4BBMBAgAiBQJXwJKWAhsDBgsJCAcDAgYVCAIJCgsE
FgIDAQIeAQIXgAAKCRD+xeKuvKsOaqIWB/4iCKYZiSgoXGVLZ3EMXnblQJvPBB4P
c5RrwprOieC3u9hdmGiZgvQDq0bZsU9iTsBqNrlEXeg7FqCskIE555IpMzVd1Kkx
PKl1KZRLWbmHuhgX+Qafygg7sAI2uzLqFp7G04az4HG4Y9sHNURWMaUE6d+aN2Nu
cRhx5icfb78iznayEIcfqlgSXuiAL93TCWHTcChwYwu7Mp3JcqfJZ+bTEw1bZMlg
R+I+qowlHtN9JJxf5Dz5HqdlWHkOqJxGfWx+SwSnSjKoa7KJBCLVydQ9Y1HkTIKT
v1fU5WNcrRVIZuAi98YueBJki8JhDM+D2s6rgJ0LZ7l3DDyT0Wlzu08WuQENBFfA
kpYBCADsGcrz0FnferYFqpYc8xqYdiP0HddVUf4Zx6AXovrk6hWSR3ejiGMfUqVl
N/a4QOjdNuBgL5X9LMEOTu9o97B7jJEGcUcQ7Ll27qMmbEwaajHF3/HsVXlQnlAE
/LsIPh90haVKqKTm9Thc4jZhxpuxsBc1FE2v24E4WLucnpyA71LHli5q8sk9ekVw
A3Iu6JswMos+Jo6OhA2YEo2VxRWS+AJhKG8XKFLrb0HsU8Sga4KiTqbasWS9iqeI
9oGOa1i4g6Lig+viv6WOWlaM+b0yK8IDzHAf03K0KaU+ZwRbuEQ056BRoDtaTw4T
I05V9c9b43pft0EftlquCWtA850zABEBAAGJAR8EGAECAAkFAlfAkpYCGwwACgkQ
/sXirryrDmoHAQf9GUPdjUISQEOJYnPjb24jQydL21+HYc01tQ7hDvGfYtR3cMzZ
Z7fovj5xyrNzP5GSY7eEywo8WiHoYNxFcUCljcjtDm143BN62S0vzw75b0nSZltX
pjJA5sUQJruuLu+E16drUxc5YaNhQooHLmso8GLRlYWbyyDAQNeXz9crxJe2HP42
6NVg2gGN1W1lHGj2j64pTkRUo6Sfkvp9rPOSRAWSBpKTn50yIkO5cXBRWTiwmsoJ
7N4unhakBD/TP9IVMPvoJ1QFe5jGMufDe1jNzFSanH/RIG924YNcxIwFJjM0Te4F
89u3LxK8yYGSSpBexff+ELmazovgI022OEOLYA==
=OkIo
-----END PGP PUBLIC KEY BLOCK-----" | apt-key add -
