# deb-meta-package/postinst
systemctl enable prex
/etc/cron.daily/auto-update-pylinkbot
