# deb-meta-package/prex.service
[Unit]
Description=Python remote execution engine server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/prex-server.py --host 0.0.0.0 --port 43000
Restart=always
User=nobody
Environment=PATH=/usr/bin:/usr/local/bin

[Install]
WantedBy=multi-user.target
